# docker/web/service/resque_scheduler/run
#!/bin/sh
cd /home/app

exec /sbin/setuser app rake environment resque:scheduler
